# deploy/ealta-frame.timer
[Unit]
Description=Advance the Inky station gallery every five minutes

[Timer]
OnBootSec=45
OnUnitActiveSec=5min
AccuracySec=5s

[Install]
WantedBy=timers.target
